# infra/systemd/templates/chat-postgresql.service.in
[Unit]
Description=Chat local PostgreSQL
After=network.target

[Service]
Type=oneshot
WorkingDirectory=@ROOT@
ExecStart=@ROOT@/scripts/infrastructure/postgresql/start.sh
ExecStop=@ROOT@/scripts/infrastructure/postgresql/stop.sh
RemainAfterExit=yes

[Install]
WantedBy=chat-infrastructure.target
